# supervisord.conf
[supervisord]
nodaemon=true

[program:redis-server]
command=/usr/bin/redis-server

[program:mariadb]
command=/bin/bash /usr/local/bin/docker-entrypoint.sh mysqld
