# supervisord.conf
[supervisord]
nodaemon=true

[program:redis-server]
command=/usr/bin/redis-server

[program:mariadb]
command=sudo -u mysql /usr/sbin/mysqld
